# examples/systemd/open-webui-speech-cache-prune.service
[Unit]
Description=Prune Open WebUI speech cache

[Service]
Type=oneshot
ExecStart=%h/open-webui-gemini-tts-proxy/scripts/prune-open-webui-speech-cache.sh
